# wetterstation.service
[Unit]
Description=Wetterstation Display
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=CHANGEME
WorkingDirectory=/home/CHANGEME/wetterstation
ExecStart=/usr/bin/python3 -m wetterstation --start
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
